# config/tribe/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
user=root

; ─── php-fpm ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ─── inotify watcher (handles startup full-index + live change watching) ──────
; watcher.sh already runs `php /app/index_files.php` on startup before watching,
; so no separate one-shot program is needed.
[program:file-watcher]
command=/bin/bash /app/watcher.sh
autostart=true
autorestart=true
priority=20
; Give php-fpm a moment to start before the watcher begins
startsecs=3
stdout_logfile=/var/log/file-watcher.log
stderr_logfile=/var/log/file-watcher-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ─── DB indexer — one-shot full backfill on startup ───────────────────────────
; Runs php /app/index_db.php once when the container starts.
; autorestart=false so it does not loop after completion.
; priority=25 so it starts after php-fpm (10) and the file watcher (20).
; After the initial run, ongoing DB → Typesense sync is handled in-process by
; Core::syncToTypesense (dual-write on every pushObject/deleteObject call).
; To force a manual re-index:  supervisorctl start db-indexer
[program:db-indexer]
command=php /app/index_db.php
autostart=true
autorestart=false
priority=25
startsecs=5
stdout_logfile=/var/log/db-indexer.log
stderr_logfile=/var/log/db-indexer-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
