# config/tribe/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
user=root

; ─── php-fpm ──────────────────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
startretries=3
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ─── DB indexer — one-shot full backfill on startup ───────────────────────────
; Runs once after php-fpm is up. autorestart=false so it exits cleanly without
; looping. priority=20 puts it before the file watcher so the DB and file
; indexes are both warm before the watcher starts its own startup check.
; Ongoing real-time sync is handled by Core::syncToTypesense on every
; pushObject / deleteObject call — this script only covers the initial backfill
; and manual re-indexes via: supervisorctl start db-indexer
[program:db-indexer]
command=php /var/www/html/config/index_db.php --skip-if-indexed
autostart=true
autorestart=false
priority=20
; Wait for php-fpm (priority 10) to be accepting connections before hitting MySQL.
startsecs=5
startretries=3
exitcodes=0,1
stdout_logfile=/var/log/db-indexer.log
stderr_logfile=/var/log/db-indexer-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ─── inotify watcher (live file change watching + conditional startup index) ──
; watcher.sh checks whether the Typesense 'files' collection is already
; populated before deciding to run a full re-index, so a container restart
; does NOT trigger an unnecessary full re-index (fixes the restart regression).
; autorestart=true is correct here — if inotifywait dies, we want it back.
; priority=30 ensures php-fpm and the DB indexer have both started first.
[program:file-watcher]
command=/bin/bash /var/www/html/config/watcher.sh
autostart=true
autorestart=true
priority=30
; Give php-fpm and the DB indexer time to settle before the watcher starts.
startsecs=5
startretries=5
stdout_logfile=/var/log/file-watcher.log
stderr_logfile=/var/log/file-watcher-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB